# data/systemd/myasus4linux-charge-limit.service
[Unit]
Description=Apply saved battery charge limit
After=multi-user.target
ConditionPathExists=/sys/class/power_supply/BAT0/charge_control_end_threshold

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for home in /home/*; do f="$home/.config/myasus4linux/settings.toml"; if [ -f "$f" ]; then v=$(grep charge_threshold "$f" | grep -o "[0-9]*"); [ -n "$v" ] && echo "$v" > /sys/class/power_supply/BAT0/charge_control_end_threshold && break; fi; done'

[Install]
WantedBy=multi-user.target
