# data/99-myasus4linux.rules
# myasus4linux: make the ASUS hardware controls writable WITHOUT a password
# prompt, but only by members of the "myasus" group -- never world-writable.
#
# These controls are range-checked by the kernel, so handing write access to a
# trusted group is safe; 0666 (world-writable) was not, since any local process
# could change the fan profile or charge limit. The app falls back to a polkit
# prompt when this rule (or the group) isn't installed.
#
# Requires the "myasus" group to exist and the user to be a member
# (see data/setup-permissions.sh, run once as root).
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="asus::kbd_backlight", RUN+="/bin/chgrp myasus /sys/class/leds/asus::kbd_backlight/brightness", RUN+="/bin/chmod 0664 /sys/class/leds/asus::kbd_backlight/brightness"
ACTION=="add", SUBSYSTEM=="platform", KERNEL=="asus-nb-wmi", RUN+="/bin/chgrp myasus /sys/devices/platform/asus-nb-wmi/throttle_thermal_policy", RUN+="/bin/chmod 0664 /sys/devices/platform/asus-nb-wmi/throttle_thermal_policy"
ACTION=="add", SUBSYSTEM=="power_supply", KERNEL=="BAT0", RUN+="/bin/chgrp myasus /sys/class/power_supply/BAT0/charge_control_end_threshold", RUN+="/bin/chmod 0664 /sys/class/power_supply/BAT0/charge_control_end_threshold"
ACTION=="add", SUBSYSTEM=="backlight", RUN+="/bin/chgrp myasus /sys/class/backlight/%k/brightness", RUN+="/bin/chmod 0664 /sys/class/backlight/%k/brightness"
