# data/myasusd.service.in
[Unit]
Description=myasus4linux privileged hardware control daemon
Documentation=https://github.com/berkmuftuoglu/myasus4linux
After=systemd-udevd.service
# Cap restarts so a persistent failure surfaces instead of crash-looping.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=dbus
BusName=io.github.berkmuftuoglu.MyAsus4Linux.Helper
ExecStart=@libexecdir@/myasusd
# Survive a transient crash so the saved charge limit still gets restored.
Restart=on-failure
RestartSec=1
# Root-owned /var/lib/myasus4linux for the persisted charge limit; systemd
# creates it and adds it to the writable paths automatically.
StateDirectory=myasus4linux
# It only ever writes a handful of fixed /sys files; lock everything else down.
# INVARIANT: every path any Op can resolve to must be under one of these roots,
# or the write hits EROFS under ProtectSystem=strict. asus-nb-wmi covers the WMI
# fan node; firmware/acpi covers the platform_profile fallback used on non-ROG
# laptops without the WMI node; power_supply covers any BATx + the charge limit;
# leds covers the keyboard backlight.
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/sys/devices/platform/asus-nb-wmi /sys/firmware/acpi /sys/class/power_supply /sys/class/leds
NoNewPrivileges=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
# It only ever talks D-Bus over a unix socket; it needs no network.
RestrictAddressFamilies=AF_UNIX
# Drop all capabilities: as root the daemon owns the sysfs files it writes, so it
# needs none. Re-verify writes still succeed if this is ever tightened further.
CapabilityBoundingSet=
AmbientCapabilities=
MemoryDenyWriteExecute=yes
# Allowlist the service syscall set; the daemon does file I/O + D-Bus + signals.
# (Verify on a real boot if you tighten this further -- a bad filter SIGSYSes.)
SystemCallFilter=@system-service
IPAddressDeny=any
RemoveIPC=yes
LockPersonality=yes
SystemCallArchitectures=native
PrivateTmp=yes
# NB: not ProtectKernelTunables -- that remounts /sys read-only and would block
# the very writes this daemon exists to make.

# Enabling the daemon (systemctl enable myasusd) makes systemd start it at boot,
# where it re-applies the saved charge limit. It stays D-Bus activatable too.
[Install]
WantedBy=multi-user.target
